# mini-unix/fs/diskA/usr/lib/crontab
0	0,12	*	*	*	/usr/ltr/SCANp > /usr/ltr/ascii
